# supervisord.conf
[supervisord]
logfile=/var/log/alt-fast/supervisord.log
pidfile=/var/log/alt-fast/supervisord.pid
childlogdir=/var/log/alt-fast
nodaemon=false

[unix_http_server]
file=/var/log/alt-fast/supervisor.sock

[supervisorctl]
serverurl=unix:///var/log/alt-fast/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:api]
command=%(here)s/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=%(here)s/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/alt-fast/api.log
stderr_logfile=/var/log/alt-fast/api.err.log

[program:news]
command=%(here)s/backend/.venv/bin/alt news collect
directory=%(here)s/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/alt-fast/news.log
stderr_logfile=/var/log/alt-fast/news.err.log

[program:market]
command=%(here)s/backend/.venv/bin/alt market collect
directory=%(here)s/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/alt-fast/market.log
stderr_logfile=/var/log/alt-fast/market.err.log

[program:dart]
command=%(here)s/backend/.venv/bin/alt dart collect
directory=%(here)s/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/alt-fast/dart.log
stderr_logfile=/var/log/alt-fast/dart.err.log

[program:trader]
command=%(here)s/backend/.venv/bin/alt trader run
directory=%(here)s/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/alt-fast/trader.log
stderr_logfile=/var/log/alt-fast/trader.err.log

[program:ws]
command=%(here)s/backend/.venv/bin/alt ws subscribe
directory=%(here)s/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/alt-fast/ws.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/alt-fast/ws.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2

[group:collectors]
programs=news,market,dart,trader,ws
